# supervisord.conf
[supervisord]
nodaemon=true

[program:smsd]
command=/usr/local/bin/smsd -t
