# supervisord.conf
[supervisord]
nodaemon=true

[program:smsd]
command=smsd -t
